# salt-work/test and stu/fileserver/ainit/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
sh /data/web/start_security.sh
echo '/data/web/logs/core/core_%e-%p-%t' >/proc/sys/kernel/core_pattern
